# launchd/com.thedoctor.start-all.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thedoctor.start-all</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/asadpreuss-dodhy/Documents/Development/the-doctor/start-all.sh</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/asadpreuss-dodhy/Documents/Development/the-doctor</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/asadpreuss-dodhy/Documents/Development/the-doctor/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/asadpreuss-dodhy/Documents/Development/the-doctor/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:__PATH_REPLACE__</string>
        <key>HOME</key>
        <string>/Users/asadpreuss-dodhy</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
